# applications/luci-app-n2n_v3/root/etc/uci-defaults/luci-n2n_v3
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@n2n_v3[-1]
	add ucitrack n2n_v3
	set ucitrack.@n2n_v3[-1].init=n2n_v3
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
